# etc/cron.weekly/certbot
#!/bin/bash

schroot -c chroot-letsencrypt -- certbot renew >> /var/log/certbot.log
